# borg-em.service
[Unit]
Description=Borg-em backup

[Service]
Type=simple
KillMode=none
ExecStart=/usr/bin/borg-em